# hotread.crontab
00 7,10,12,15,18,20 * * * python /root/hotread/social/ggweibo.py >>/root/log/ggweibo.log 2>&1
00 * * * * python /root/hotread/feed/feed_mutli.py >>/root/log/feed_mutli.log 2>&1
00 7,10,12,15,18,20 * * * /usr/bin/python /root/toy/weibo_hot_news.py >>/root/log/weibo_hot_news.log 2>&1
